# Yang-T_A9/crontab
0 3 * * 0 /Users/tiger/mini-assignments-tigery2005/Yang-T_A9/backup.sh
0 2 * * * rm -rf /Users/tiger/mini-assignments-tigery2005/Yang-T_A9/temp/*
30 8 * * * df -h | mail -s "Daily Disk Usage Report" tigeryang2005@gmail.com
